# launchd/com.skcapstone.skcomm-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skcapstone.skcomm-daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>${HOME}/.skenv/bin/skcomm</string>
    <string>daemon</string>
    <string>--all-agents</string>
    <string>--interval</string>
    <string>5</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>${HOME}/.skenv/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>${HOME}/.skcomm/daemon.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>${HOME}/.skcomm/daemon.stderr.log</string>
</dict>
</plist>
